# rc.local
#!/bin/sh -e
python3 /home/sergio/Desktop/netTablesTest.py & 
^this should be location of python file(delete this line)
